# systemd/openclaw-telegram-api-proxy.service.example
[Unit]
Description=OpenClaw Telegram Bot API proxy с local/cloud fallback
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=%h/.openclaw
Environment=LISTEN_HOST=127.0.0.1
Environment=PORT=8082
Environment=LOCAL_API_ROOT=http://127.0.0.1:8081
Environment=CLOUD_API_ROOT=https://api.telegram.org
Environment=ENABLE_CLOUD_FALLBACK=1
Environment=TELEGRAM_OFFSET_DIR=telegram
Environment=CLOUD_FILE_FALLBACK_MAX_BYTES=20971520
Environment=BUFFER_LIMIT_BYTES=8388608
Environment=LOCAL_HEALTH_TTL_MS=5000
Environment=LOCAL_UNHEALTHY_COOLDOWN_MS=5000
Environment=LOCAL_HEALTH_TIMEOUT_MS=2000
Environment=UPSTREAM_TIMEOUT_MS=130000
ExecStart=/usr/bin/node %h/.openclaw/bin/telegram-bot-api-proxy.mjs
Restart=always
RestartSec=3
TimeoutStartSec=20
TimeoutStopSec=10
StandardOutput=append:%h/.openclaw/logs/telegram-bot-api-proxy.log
StandardError=append:%h/.openclaw/logs/telegram-bot-api-proxy.log

[Install]
WantedBy=default.target
